# SamsungGalaxyNote/samsung-galaxy-book-linux-fixes-main/webcam-fix/v4l2-relayd-watchdog.service
[Unit]
Description=v4l2-relayd health watchdog (blank frame detection)
ConditionPathExists=/dev/video0

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/v4l2-relayd-watchdog.sh
