# SamsungGalaxyNote/samsung-galaxy-book-linux-fixes-main/webcam-fix/v4l2-relayd-watchdog.timer
[Unit]
Description=Run v4l2-relayd health watchdog every 3 minutes

[Timer]
OnActiveSec=3min
OnUnitActiveSec=3min
RandomizedDelaySec=10

[Install]
WantedBy=timers.target
