# SamsungGalaxyNote/samsung-galaxy-book-linux-fixes-main/speaker-fix/max98390-hda-check-upstream.service
[Unit]
Description=Check if native MAX98390 HDA support has landed upstream
After=multi-user.target
# Only run if the DKMS package is still installed
ConditionPathExists=/usr/src/max98390-hda-1.0/dkms.conf

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/max98390-hda-check-upstream.sh

[Install]
WantedBy=multi-user.target
