# SamsungGalaxyNote/samsung-galaxy-book-linux-fixes-main/webcam-fix-book5/ipu-bridge-check-upstream.service
[Unit]
Description=Check if kernel ipu-bridge has Samsung camera rotation fix
After=multi-user.target
# Only run if the DKMS package is still installed
ConditionPathExists=/usr/src/ipu-bridge-fix-1.0/dkms.conf

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/ipu-bridge-check-upstream.sh

[Install]
WantedBy=multi-user.target
